# etc/udev/rules.d/70-displays.rules
KERNEL=="card0", SUBSYSTEM=="drm", ACTION=="change", RUN+="/usr/bin/su bobes -c /usr/local/bin/setup_displays"
